# debian-package/cockpit-domain-controller_1.0.73-1/usr/share/cockpit/domain-controller/domain-service-orchestrator.service
[Unit]
Description=Domain Service Orchestrator
Documentation=man:domain-service-orchestrator(8)
After=samba-ad-dc.service
Wants=samba-ad-dc.service
PartOf=samba-ad-dc.service

[Service]
Type=oneshot
ExecStart=/home/dguedry/Documents/ad-server/cockpit-domain-controller/domain-service-orchestrator.sh --orchestrate
User=root
Group=root
StandardOutput=journal
StandardError=journal
TimeoutStartSec=300
RemainAfterExit=no

# Security settings
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/etc/chrony /etc/dhcp /var/lib/samba/sysvol /var/run /var/log
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target